# Scheduler/crontab
# crontab 파일을 열어 편집
crontab -e

# 다음 라인을 추가하여 서버 부팅 시 스크립트 자동 실행
# 수정 필요 도커환경인 점 반영
0 1 * * * /usr/bin/python3 /home/fuser01/scripts/FinancialSatementsCall.py >> /home/fuser01/scripts/cron.log 2>&1
